# streaming/supervisord.ini
[program:osrs-client]
command=/app/start.sh
environment=DISPLAY=":0",HOME="/root",USER="root"
autorestart=false
priority=200
stdout_logfile=/var/log/osrs-client.log
stderr_logfile=/var/log/osrs-client.log
